# gaseous_server/rootfs/etc/services.d/gaseous/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Gaseous Server
# ==============================================================================

declare DB_HOST
declare DB_USER
declare DB_PASS
declare DB_NAME
declare IGDB_CLIENT_ID
declare IGDB_CLIENT_SECRET
declare PORT

bashio::log.info "Starting Gaseous Server..."

# Get configuration
DB_HOST=$(bashio::config 'database_host')
DB_USER=$(bashio::config 'database_user')
DB_PASS=$(bashio::config 'database_password')
DB_NAME=$(bashio::config 'database_name')
IGDB_CLIENT_ID=$(bashio::config 'igdb_client_id')
IGDB_CLIENT_SECRET=$(bashio::config 'igdb_client_secret')
PORT=$(bashio::config 'port')

bashio::log.info "Configuration loaded"
bashio::log.info "Database: ${DB_HOST}/${DB_NAME}"
bashio::log.info "Port: ${PORT}"

# Wait for database
bashio::log.info "Testing database connection..."
RETRY=0
until mysql -h "${DB_HOST}" -u "${DB_USER}" -p"${DB_PASS}" -e "SELECT 1;" &>/dev/null || [ ${RETRY} -eq 30 ]; do
    bashio::log.info "Waiting for database... (${RETRY}/30)"
    RETRY=$((RETRY+1))
    sleep 2
done

if [ ${RETRY} -eq 30 ]; then
    bashio::log.warning "Could not connect to database, starting anyway..."
fi

# Create directories
mkdir -p /data/roms /data/config

# Start Gaseous
bashio::log.info "Starting Gaseous Server..."
cd /opt/gaseous || exit 1

exec java -jar gaseous-server.jar \
    --database.host="${DB_HOST}" \
    --database.user="${DB_USER}" \
    --database.password="${DB_PASS}" \
    --database.name="${DB_NAME}" \
    --igdb.client-id="${IGDB_CLIENT_ID}" \
    --igdb.client-secret="${IGDB_CLIENT_SECRET}" \
    --server.port="${PORT}" \
    --data.dir="/data/roms"
