# blinko-addon/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PostgreSQL service
# ==============================================================================

declare postgres_password
declare timezone

bashio::log.info "Starting PostgreSQL..."

# Get configuration
postgres_password=$(bashio::config 'postgres_password')
timezone=$(bashio::config 'timezone')

# Initialize PostgreSQL data directory if it doesn't exist
if [ ! -d "/data/postgres/data" ]; then
    bashio::log.info "Initializing PostgreSQL database..."
    mkdir -p /data/postgres/data
    chown -R postgres:postgres /data/postgres
    
    su - postgres -c "initdb -D /data/postgres/data"
    
    # Configure PostgreSQL
    echo "host all all 0.0.0.0/0 md5" >> /data/postgres/data/pg_hba.conf
    echo "listen_addresses = '*'" >> /data/postgres/data/postgresql.conf
    
    # Start PostgreSQL temporarily to set password
    su - postgres -c "pg_ctl -D /data/postgres/data -o \"-c listen_addresses='localhost'\" -w start"
    
    # Create database and user
    psql -U postgres -c "ALTER USER postgres WITH PASSWORD '${postgres_password}';"
    psql -U postgres -c "CREATE DATABASE postgres;"
    
    su - postgres -c "pg_ctl -D /data/postgres/data -m fast -w stop"
fi

# Set timezone
export TZ="${timezone}"

# Start PostgreSQL
bashio::log.info "Starting PostgreSQL server..."
exec su - postgres -c "postgres -D /data/postgres/data"
